# container-files/etc/supervisor.d/memcached.conf
[program:memcached]
command = /usr/bin/memcached -u www -p 11211 -r -m 256 -c 1024 -P /data/run/memcached/memcached
user = www
autorestart = false
stderr_logfile = NONE
stdout_logfile = NONE
